# speedtest.timer
[Unit]
Description=Run a speedtest every 15 minutes

[Timer]
OnCalendar=*:0/15
# RandomizedDelaySec=30

[Install]
WantedBy=timers.target
